# infra/macmini/launchd/com.tokensoff.colima.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tokensoff.colima</string>
  <key>ProgramArguments</key>
  <array>
    <string>__PROJECT_DIR__/infra/macmini/run-colima.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>StandardOutPath</key>
  <string>__LOG_DIR__/colima.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/colima.error.log</string>
</dict>
</plist>

